# recipes-graphics/x11-common/files/xserver-nodm-local.service.in
[Unit]
Description=Xserver startup without a display manager

[Service]
EnvironmentFile=/etc/default/xserver-nodm
Environment="LD_LIBRARY_PATH=/usr/local/lib/"
User=@USER@
ExecStart=/etc/xserver-nodm/Xserver
Restart=always

[Install]
Alias=display-manager.service
